# services/redis/run
#!/bin/sh
exec s6-notifyoncheck -d -s 250 s6-setuidgid redis redis-server --port 6379 --bind 127.0.0.1 

